# deploy/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for FCR Gate
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=0

[Service]
Type=notify
TimeoutStartSec=0
ExecStart=/data/fcr-gate/bin/cloudflared tunnel --no-autoupdate run --token-file /data/fcr-gate/secrets/cloudflare-tunnel-token
Restart=always
RestartSec=5s
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
